# netmonitor-feed-update.service
################################################################################
# ⚠️  WARNING: THIS IS A TEMPLATE FILE - DO NOT COPY DIRECTLY TO /etc/systemd/system/
################################################################################
#
# See netmonitor.service for installation instructions.
# The install_complete.sh script will generate this file with correct paths.
#
################################################################################

[Unit]
Description=Network Monitor - Threat Feed Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/venv/bin/python3 __INSTALL_DIR__/update_feeds.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
