# netmonitor-sensor.service
# NOTE: This is a TEMPLATE file for remote sensors. Do not use directly!
#
# Installation steps for Nano Pi / remote sensor:
#
# 1. Install NetMonitor in /opt/netmonitor (via install_complete.sh)
#
# 2. Create sensor.conf from template:
#    cd /opt/netmonitor
#    cp sensor.conf.template sensor.conf
#    nano sensor.conf
#    - Set SOC_SERVER_URL
#    - Set SENSOR_ID (unique per sensor)
#    - Set SENSOR_LOCATION
#    - Set SENSOR_TOKEN (generate via setup_sensor_auth.py on SOC server)
#    - Set INTERFACE (eth0, eth1, etc.)
#
# 3. Copy service file and enable:
#    sudo cp netmonitor-sensor.service /etc/systemd/system/
#    sudo systemctl daemon-reload
#    sudo systemctl enable netmonitor-sensor
#    sudo systemctl start netmonitor-sensor
#
# 4. Check status:
#    sudo systemctl status netmonitor-sensor
#    sudo journalctl -u netmonitor-sensor -f

[Unit]
Description=NetMonitor Remote Sensor Client
Documentation=https://github.com/yourusername/netmonitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netmonitor
ExecStart=/opt/netmonitor/venv/bin/python3 /opt/netmonitor/sensor_client.py -c /opt/netmonitor/sensor.conf

# NOTE: Config is now in sensor.conf (bash KEY=value format)
# Edit /opt/netmonitor/sensor.conf to set:
#   SOC_SERVER_URL, SENSOR_ID, SENSOR_LOCATION, SENSOR_TOKEN, INTERFACE

# Restart on failure
Restart=on-failure
RestartSec=10

# Security
PrivateTmp=yes
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netmonitor-sensor

[Install]
WantedBy=multi-user.target
